# apps/master/scripts/autmzr-command-backup.timer
[Unit]
Description=autmzr-command DB backup daily at 03:30 UTC

[Timer]
OnCalendar=*-*-* 03:30:00 UTC
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
