# apps/master/scripts/pocket-claude-backup.service
[Unit]
Description=pocket-claude DB backup (pg_dump)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
EnvironmentFile=/opt/pocket-claude/.env
Environment=BACKUP_DIR=/var/backups/pocket-claude
Environment=RETENTION_DAILY=14
Environment=RETENTION_WEEKLY=8
Environment=BACKUP_LOG=/var/log/pocket-claude-backup.log
ExecStart=/opt/pocket-claude/apps/master/scripts/backup.sh
StandardOutput=append:/var/log/pocket-claude-backup.log
StandardError=append:/var/log/pocket-claude-backup.log
Nice=15
IOSchedulingClass=idle
